# security/sentinel/runit/karya-sentinel/run
#!/usr/bin/env bash
# runit servisi icin Karya Sentinel
exec /usr/lib/karya/scripts/karya-sentinel.sh start
