# hardware/speed/services/karya-speed.path
[Unit]
Description=Karya SPEED - Profile Watcher

[Path]
PathModified=/etc/karya/speed/current.conf

[Install]
WantedBy=multi-user.target
